# applications/luci-app-xupnpd/root/etc/uci-defaults/luci-xupnpd
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@xupnpd[-1]
	add ucitrack xupnpd
	set ucitrack.@xupnpd[-1].init=xupnpd
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
