# applications/luci-app-cd8021x/root/etc/uci-defaults/luci-cd8021x
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@cd8021x[-1]
	add ucitrack cd8021x
	set ucitrack.@cd8021x[-1].init=cd8021x
EOF

rm -rf /tmp/luci-indexcache*
exit 0
