# modules/luci-oldapi/root/etc/uci-defaults/099-oldapi-copy
#!/bin/sh

cp -a /etc/oldapi/* /

exit 0
