# applications/luci-app-minieap/root/etc/uci-defaults/luci-minieap
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@minieap[-1]
	add ucitrack minieap
	set ucitrack.@minieap[-1].init=minieap
	commit ucitrack
EOF

rm -rf /tmp/luci-*
exit 0
