# applications/luci-app-onekeyap/root/etc/hotplug.d/iface/99-onekeyap
#!/bin/sh

[ "$ACTION" = "ifup" ] || [ "$ACTION" = "ifupdate" ] || exit 0
[ "$INTERFACE" = "lan" ] || exit 0

/usr/libexec/onekeyap refresh-fallback >/dev/null 2>&1
