# applications/luci-app-kodexplorer/root/etc/uci-defaults/luci-app-kodexplorer
#!/bin/sh

/etc/init.d/php8-fpm disable && /etc/init.d/php8-fpm stop
rm -rf /tmp/luci-*cache
exit 0
